# osclass/apache2.ini
[program:apache2]
command=/usr/sbin/httpd -D FOREGROUND
autostart=true
autorestart=true
stderr_logfile = /var/log/supervisord/apache2-stderr.log
stdout_logfile = /var/log/supervisord/apache2-stdout.log

